# glam-package/usr/src/glam/tareasUsuarios/subScripts/respaldoXsession/demonInOut.service
[Unit]
Description=Detecta inicios y terminos de sesión
After=network.target

[Service]
ExecStart=/etc/init.d/demonScript.sh

[Install]
WantedBy=default.target